# tests/fixtures/schedule/android-watcher.service
[Unit]
Description=android-watcher scheduled run

[Service]
Type=oneshot
ExecStart=/usr/bin/android-watcher run
